# packaging/sysknife-daemon.service
[Unit]
Description=LACS privileged daemon
Documentation=https://github.com/lacs-foundation/sysknife
After=network.target

[Service]
Type=simple
User=sysknife
Group=sysknife

# Production socket and database paths (overridable via drop-in).
Environment="SYSKNIFE_LISTEN_URI=unix:///run/sysknife/daemon.sock"
Environment="SYSKNIFE_DATABASE_PATH=/var/lib/sysknife/daemon.sqlite"

ExecStart=/usr/local/bin/sysknife-daemon
Restart=on-failure
RestartSec=5s

# Harden the service:
# - NoNewPrivileges prevents the daemon itself from gaining additional
#   privileges after launch; polkit rules grant per-action elevation.
# - ProtectHome and ProtectSystem prevent reads outside /var/lib/sysknife.
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/sysknife /run/sysknife
RuntimeDirectory=sysknife
StateDirectory=sysknife

[Install]
WantedBy=multi-user.target
